# rpm_extracts/openvswitch/systemd/usr_lib_systemd_system_ovs-delete-invalid-ports.service
[Unit]
Description=Open vSwitch Delete Invalid Ports
After=ovs-vswitchd.service
Before=network.target network.service

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/usr/share/openvswitch/scripts/ovs-ctl delete-invalid-ports
