# rpm_extracts/rpcbind/systemd/usr_lib_systemd_system_rpcbind.socket
[Unit]
Description=RPCbind Server Activation Socket
DefaultDependencies=no
Wants=rpcbind.target
Before=rpcbind.target

[Socket]
ListenStream=/run/rpcbind.sock

[Install]
WantedBy=sockets.target
