# rpm_extracts/libvirt-daemon-driver-nwfilter/systemd/usr_lib_systemd_system_virtnwfilterd-ro.socket
[Unit]
Description=Libvirt nwfilter local read-only socket
Before=virtnwfilterd.service
BindsTo=virtnwfilterd.socket
After=virtnwfilterd.socket
Conflicts=libvirtd.socket libvirtd-ro.socket libvirtd-admin.socket libvirtd-tcp.socket libvirtd-tls.socket

[Socket]
# The directory must match the /etc/libvirt/virtnwfilterd.conf unix_sock_dir setting
# when using systemd version < 227
ListenStream=/var/run/libvirt/virtnwfilterd-sock-ro
Service=virtnwfilterd.service
SocketMode=0666

[Install]
WantedBy=sockets.target
